# rootfs/etc/services.d/nestor-bridge/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start Nestor Bridge
# ==============================================================================
bashio::log.info "Starting Nestor Bridge..."

exec dotnet /app/NestorBridge.dll
